# deploy/mem.supervisor.conf
# mem — supervisor program template (manual deployment).
#
# This is the program block `scripts/install.sh` generates for supervisor;
# provided here so you can deploy by hand. Steps:
#   1. Put the binary at /usr/local/bin/mem (build: `cargo build --release`).
#   2. Create the service user + dirs:
#        sudo useradd --system --home-dir /var/lib/mem --shell /usr/sbin/nologin mem
#        sudo install -d -m 0750 -o mem -g mem /var/lib/mem /var/lib/mem/hf-cache /var/log/mem
#   3. Drop deploy/mem.config.env.example at /var/lib/mem/config.env (edit it).
#   4. Copy this file to /etc/supervisor/conf.d/mem.conf, then:
#        sudo supervisorctl reread && sudo supervisorctl update
#        sudo supervisorctl status mem ; sudo supervisorctl tail -f mem
#
# supervisor has no EnvironmentFile= — the command sources config.env in a
# sh -c wrapper so the env file is the single source of truth (shared with
# the systemd unit). DuckDB is single-writer — one program per MEM_DB_PATH.

[program:mem]
command=/bin/sh -c 'set -a; . /var/lib/mem/config.env; set +a; exec /usr/local/bin/mem serve'
directory=/var/lib/mem
user=mem
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/var/log/mem/mem.out.log
stderr_logfile=/var/log/mem/mem.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
